# Lampi/pkg/lampi/etc/supervisor/conf.d/lamp_ui.conf
[program:lamp_ui]
command=/opt/lampi/lamp_ui.py
directory=/opt/lampi
environment=HOME="/home/pi"
priority=400
user=pi
autostart=true
autorestart=true
